# data/eos-metrics-instrumentation.service.in
[Unit]
Description=EndlessOS Metrics Instrumentation
Requires=eos-metrics-event-recorder.service
After=eos-metrics-event-recorder.service

[Service]
Type=simple
ExecStart=@libexecdir@/eos-metrics-instrumentation
User=metrics

[Install]
WantedBy=multi-user.target
