# data/eos-metrics-instrumentation.path.in
[Unit]
Description=Endless OS Metrics Instrumentation trigger
Wants=geoclue.service
After=eos-metrics-event-recorder.service

[Path]
DirectoryNotEmpty=/home

[Install]
WantedBy=multi-user.target
