# data/eos-metrics-instrumentation.service.in
[Unit]
Description=Endless OS Metrics Instrumentation
Requires=eos-metrics-event-recorder.service
Wants=geoclue.service
After=eos-metrics-event-recorder.service
ConditionDirectoryNotEmpty=/home

[Service]
Type=simple
ExecStart=@libexecdir@/eos-metrics-instrumentation
